# deployment/sizif/etc/cron.d/logrotate
0 0 * * * root /usr/sbin/logrotate -f /etc/logrotate.d/sizif
